# templates/radar-output.service.epp
<%- | String  $user, String $docker_repo_dir, String $volume_dir,
| -%>
[Unit]
Description=RADAR-Docker output service
After=radar-docker.service
Requires=docker.service

[Service]

WorkingDirectory=/home/<%= $user -%>/radar-docker<%= $docker_repo_dir %>

ExecStart=/home/<%= $user -%>/radar-docker<%= $docker_repo_dir %>/bin/hdfs-restructure /topicAndroidNew <%= $volume_dir %>/output

NotifyAccess=all

