# templates/radar-health.service.epp
<%- | String  $user, String $docker_repo_dir,
| -%>
[Unit]
Description=RADAR-Docker check-health service
After=radar-docker.service
Requires=docker.service

[Service]

WorkingDirectory=/home/<%= $user -%>/radar-docker<%= $docker_repo_dir %>

ExecStart=/home/<%= $user -%>/radar-docker<%= $docker_repo_dir %>/bin/healthcheck

NotifyAccess=all

